# deploy/systemd/trading-bot-db-backup.timer
[Unit]
Description=Daily TimescaleDB backup for AI Trading Bot

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
