# deploy/systemd/trading-bot-paper-live.service
[Unit]
Description=AI Trading Bot live paper (WS 5m + model)
After=network-online.target docker.service trading-bot-ingest.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-trading-bot
Environment=PATH=/opt/ai-trading-bot/.venv/bin:/usr/local/bin:/usr/bin
EnvironmentFile=-/root/.config/trading-bot/binance_testnet.env
ExecStart=/opt/ai-trading-bot/.venv/bin/trading-bot --metrics-port 9110 paper-live --seconds 0 --symbol BTCUSDT --interval 5m --model data/models/xgb_btc_5m.joblib --cash 1000
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
